# docker/backup/backup/backup.cron
# run the backup every day at 02:00 UTC
0 2 * * * /opt/backup/backup.sh
